# misc/ec2/testnet/scion/systemd/chrony.service
[Unit]
Description=chrony
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/chrony/chronyd.pid
ExecStartPre=timedatectl set-ntp false
ExecStart=/home/ec2-user/chrony-4.6.1/sbin/chronyd -4 -f /home/ec2-user/testnet/scion/chrony.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
