# misc/ec2/testnet/systemd/scion-dispatcher@.service
[Unit]
Description=SCION Dispatcher
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=rm -rf /run/shm/dispatcher
WorkingDirectory=/home/ec2-user/testnet
ExecStart=/home/ec2-user/scion/bin/dispatcher --config gen/%i/dispatcher.toml
LimitNOFILE=4096
Restart=on-failure

[Install]
WantedBy=multi-user.target
