# misc/ec2/testnet/ip/chrony@.service
[Unit]
Description=chrony
After=network-online.target
Wants=network-online.target

[Service]
Type=forking
PIDFile=/run/chrony/chronyd.pid
ExecStartPre=rm -rf /var/run/chrony
ExecStartPre=timedatectl set-ntp false
ExecStart=/home/ec2-user/chrony-4.4/sbin/chronyd -4 -f /home/ec2-user/testnet/ip/chrony_%i.conf
Restart=on-failure

[Install]
WantedBy=multi-user.target
